# deploy/launchd/com.ilyin-stroy.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ilyin-stroy.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/server/projects/ilyin_stroy/scripts/telegram_bot.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>NetworkState</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/server/projects/ilyin_stroy</string>

    <key>StandardOutPath</key>
    <string>/Users/server/projects/ilyin_stroy/telegram_bot_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/server/projects/ilyin_stroy/telegram_bot_stderr.log</string>
</dict>
</plist>
